# infra/cron-vm/units/hamafx-docker-prune.service
[Unit]
Description=HamaFX-Ai weekly Docker image prune (reclaim disk space)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
NoNewPrivileges=true
PrivateTmp=true
ExecStart=/usr/bin/docker image prune -f
ExecStartPost=/bin/sh -c 'test -z "$$HC_DISK_CHECK_UUID" || /usr/bin/curl -fsS -m 5 "https://hc-ping.com/$$HC_DISK_CHECK_UUID" || true'

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-docker-prune

StartLimitIntervalSec=3600
StartLimitBurst=3
