# infra/cron-vm/units/hamafx-docker-prune.service
[Unit]
Description=HamaFX-Ai weekly Docker image prune (reclaim disk space)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
NoNewPrivileges=true
PrivateTmp=true
ExecStart=/usr/bin/docker image prune -f

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-docker-prune

StartLimitIntervalSec=3600
StartLimitBurst=3
