# infra/cron-vm/units/hamafx-backup-journal.service
[Unit]
Description=HamaFX-Ai nightly journal_entries JSON export → GCS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hamafx
Group=hamafx
WorkingDirectory=/opt/hamafx
EnvironmentFile=/opt/hamafx/.env
ExecStart=/opt/hamafx/scripts/backup-journal.sh

MemoryMax=256M
TimeoutStartSec=300

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hamafx /tmp

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-backup-journal
