# infra/cron-vm/units/hamafx-job-snapshots.timer
[Unit]
Description=Daily snapshots at 00:05 UTC

[Timer]
OnCalendar=*-*-* 00:05:00 UTC
RandomizedDelaySec=60
Persistent=true
Unit=hamafx-job-snapshots.service

[Install]
WantedBy=timers.target
