# infra/cron-vm/units/hamafx-job-snapshots.service
[Unit]
Description=HamaFX-Ai job — daily snapshots + candles_1m prune
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hamafx
Group=hamafx
WorkingDirectory=/opt/hamafx/app/apps/worker
EnvironmentFile=/opt/hamafx/.env
# Phase 2 hardening §4 — pick worker-flavoured pool defaults in @hamafx/db.
Environment="HAMAFX_RUNTIME=worker"
# Phase 3 hardening §12 — force JSON logger output under journald.
Environment="NODE_ENV=production"
ExecStart=/usr/bin/node /opt/hamafx/app/apps/worker/dist/runner/cli.js snapshots

MemoryMax=512M
TimeoutStartSec=300
StartLimitIntervalSec=3600
StartLimitBurst=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hamafx

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-job-snapshots
