# infra/cron-vm/units/hamafx-update.timer
[Unit]
Description=Self-update every 5 minutes
# Run a few minutes after boot so a fresh VM picks up changes promptly.

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30
Unit=hamafx-update.service

[Install]
WantedBy=timers.target
