# infra/cron-vm/units/hamafx-job-fred-actuals.timer
[Unit]
Description=FRED actuals backfill — 01:30 UTC daily

[Timer]
OnCalendar=*-*-* 01:30:00 UTC
RandomizedDelaySec=60
Persistent=true
Unit=hamafx-job-fred-actuals.service

[Install]
WantedBy=timers.target
