# infra/cron-vm/units/hamafx-job-weekly-review.service
[Unit]
Description=HamaFX-Ai job — Sunday weekly journal review
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hamafx
Group=hamafx
WorkingDirectory=/opt/hamafx/app/apps/worker
EnvironmentFile=/opt/hamafx/.env
ExecStart=/usr/bin/node /opt/hamafx/app/apps/worker/dist/runner/cli.js weekly-review

MemoryMax=512M
TimeoutStartSec=300

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hamafx

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-job-weekly-review
