# infra/cron-vm/units/hamafx-job-cot.timer
[Unit]
Description=CFTC CoT weekly ingestion — Friday 22:00 UTC

[Timer]
OnCalendar=Fri *-*-* 22:00:00 UTC
RandomizedDelaySec=60
Persistent=true
Unit=hamafx-job-cot.service

[Install]
WantedBy=timers.target
