# infra/cron-vm/units/hamafx-job-embedding-backfill.timer
[Unit]
Description=Schedule embedding-backfill every 6 hours
# We do not Persistent=true — if the VM is down at the trigger we don't
# need to catch up on missed runs (the next run will absorb the backlog).

[Timer]
OnCalendar=0/6:00:00 UTC
RandomizedDelaySec=60
Unit=hamafx-job-embedding-backfill.service

[Install]
WantedBy=timers.target
