# infra/cron-vm/units/hamafx-update.timer
[Unit]
Description=Self-update every 5 minutes (git pull + Docker build)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30
Unit=hamafx-update.service

[Install]
WantedBy=timers.target
