# infra/cron-vm/units/hamafx-tenant-export.timer
[Unit]
Description=Weekly per-tenant export rehearsal — Sunday 05:00 UTC

[Timer]
OnCalendar=Sun *-*-* 05:00:00 UTC
RandomizedDelaySec=120
Persistent=true
Unit=hamafx-tenant-export.service

[Install]
WantedBy=timers.target
