# infra/cron-vm/units/hamafx-job-fred-actuals.service
[Unit]
Description=HamaFX-Ai job — FRED actuals backfill
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hamafx
Group=hamafx
WorkingDirectory=/opt/hamafx/app/apps/worker
EnvironmentFile=/opt/hamafx/.env
ExecStartPre=/bin/sh -c 'echo "DEPLOYED_SHA=$(cat /opt/hamafx/.deployed-sha 2>/dev/null || echo unknown)" > /run/hamafx-job-fred-actuals.env'
EnvironmentFile=-/run/hamafx-job-fred-actuals.env
ExecStart=/usr/bin/node /opt/hamafx/app/apps/worker/dist/runner/cli.js fred-actuals

MemoryMax=512M
TimeoutStartSec=120

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hamafx /run/hamafx-job-fred-actuals.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-job-fred-actuals
