# infra/cron-vm/units/hamafx-light-warm-cache.service
[Unit]
Description=HamaFX-Ai light cron — POST /api/cron/warm-cache on Vercel
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hamafx
Group=hamafx
EnvironmentFile=/opt/hamafx/.env
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hamafx
# PR-15: Warm both the market cache (/api/cron/warm-cache) and the chat endpoint
# (/api/health/public) to reduce cold starts for infrequent users.
# The chat endpoint imports 32 AI tools + multi-agent orchestration;
# keeping it warm avoids 2-5s cold start penalty on first turn.
ExecStart=/usr/bin/curl -fsS -m 30 \
  -H "Authorization: Bearer ${CRON_SECRET}" \
  ${PRODUCTION_URL}/api/cron/warm-cache
# PR-15: Lightweight poke to /api/health/public keeps the chat route's
# Vercel function warm. No auth needed, minimal DB query.
# '-' prefix prevents the unit from failing if Vercel is unreachable.
ExecStartPost=-/usr/bin/curl -fsS -m 10 ${PRODUCTION_URL}/api/health/public
ExecStartPost=/bin/sh -c 'test -z "$HC_LIGHT_WARM_CACHE_UUID" || /usr/bin/curl -fsS -m 5 https://hc-ping.com/$HC_LIGHT_WARM_CACHE_UUID'

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-light-warm-cache
