# infra/cron-vm/units/hamafx-worker.service
[Unit]
Description=HamaFX-Ai worker (BiQuote SignalR consumer + heavy job runner)
Documentation=https://github.com/HamaFx/HamaFX-Ai
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
User=hamafx
Group=hamafx
WorkingDirectory=/opt/hamafx/app
EnvironmentFile=/opt/hamafx/.env
# Phase 2 hardening §4 — tells `@hamafx/db.getDb()` to pick the
# worker-flavoured pool defaults (max=3) instead of the web defaults
# (max=5). Override with WORKER_DB_POOL_MAX in /opt/hamafx/.env if
# the worker grows new high-concurrency jobs.
Environment="HAMAFX_RUNTIME=worker"
# Phase 3 hardening §12 — `EnvironmentFile` does not set NODE_ENV; the
# logger's pretty-print path was firing under journald, producing
# lines that lost their structured fields when piped through Sentry.
Environment="NODE_ENV=production"

ExecStart=/usr/bin/node /opt/hamafx/app/apps/worker/dist/index.js

# Resilience
Restart=always
RestartSec=5
# Cgroup ceiling for the always-on consumer. The 4 GB host (e2-medium) keeps
# the remaining ~2.5 GB free for transient one-shot jobs (embedding-backfill,
# briefings) which run as their own systemd units in their own cgroups.
MemoryMax=1.5G
# Phase 2 hardening §1 — `Type=notify` plus a 120-second watchdog. The
# worker's tick handler calls sd_notify(WATCHDOG=1) (throttled to once
# per 30 s) so a healthy connection refreshes the timer four times per
# window. A genuine hang or a SignalR connection that never recovers
# triggers a kill + Restart=always cycle.
WatchdogSec=120
# Use `WatchdogSignal=SIGTERM` so the worker has a chance to flush its
# tick buffer before being forcibly killed. The default is SIGABRT,
# which would skip the graceful-shutdown path entirely.
WatchdogSignal=SIGTERM
TimeoutStartSec=60
TimeoutStopSec=30
StartLimitIntervalSec=300
StartLimitBurst=10

# Hardening (non-exhaustive — keep simple for personal-mode)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hamafx

# Logging — journald handles structured JSON natively
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-worker

[Install]
WantedBy=multi-user.target
