# infra/cron-vm/units/hamafx-update.service
[Unit]
Description=HamaFX-Ai self-update (git pull + Docker build + container restart)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=hamafx
Group=hamafx
WorkingDirectory=/opt/hamafx
EnvironmentFile=/opt/hamafx/.env
ExecStart=/opt/hamafx/scripts/docker-update.sh
TimeoutStartSec=600

# Docker needs NoNewPrivileges=false (hamafx is in docker group)
NoNewPrivileges=false

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-update
