# infra/cron-vm/units/hamafx-job-embedding-backfill.service
[Unit]
Description=HamaFX-Ai job — embedding-backfill
Documentation=https://github.com/HamaFx/HamaFX-Ai/blob/main/docs/superpowers/specs/2026-05-27-phase-8-backend-reliability-design.md
# We don't depend on the always-on worker — the job runs in its own
# Node process and connects to Postgres / the AI Gateway directly.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hamafx
Group=hamafx
WorkingDirectory=/opt/hamafx/app/apps/worker
EnvironmentFile=/opt/hamafx/.env
# Phase 2 hardening §4 — pick worker-flavoured pool defaults in @hamafx/db.
Environment="HAMAFX_RUNTIME=worker"
# Phase 3 hardening §12 — force JSON logger output under journald.
Environment="NODE_ENV=production"
ExecStart=/usr/bin/node /opt/hamafx/app/apps/worker/dist/runner/cli.js embedding-backfill

# Per-job memory cap. Embedding-backfill batches 32 articles at ~1.5 KB
# each and is the heaviest job; 1 GB is plenty.
MemoryMax=1G
TimeoutStartSec=600
StartLimitIntervalSec=3600
StartLimitBurst=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hamafx

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hamafx-job-embedding-backfill
