# tools/mt5/mt5-headless.service
[Unit]
Description=MetaTrader 5 Headless Client (Wine + Xvfb)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
Environment=WINEPREFIX=/home/ubuntu/.wine
Environment=DISPLAY=:1
ExecStartPre=-/usr/bin/Xvfb :1 -screen 0 1024x768x16 &
ExecStart=/usr/bin/wine /home/ubuntu/.mt5/drive_c/Program\ Files/MetaTrader\ 5/terminal.exe /config:startup.ini
Restart=always
RestartSec=10
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
